# lb/supervisord-haproxy.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=service haproxy start
numprocs=1
autostart=true
autorestart=true

[program:serf]
command=/start-serf.sh
